# hids.service
[Unit]
Description=HIDS Monitoring Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/your/hids_script.py
WorkingDirectory=/path/to/your/script
Restart=always
User=your_username

[Install]
WantedBy=multi-user.target